# deploy/systemd/skimmer-workflow.service
[Unit]
Description=Skimmer YouTube and profile collection workflow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orangepi
WorkingDirectory=/home/orangepi/code_projects/skimmer
Environment=YOUTUBE_HEADLESS=false
Environment=VIDIQ_HEADLESS=true
Environment=SOCIALBLADE_HEADLESS=false
Environment=SOCIALBLADE_CHANNEL_DELAY_SECONDS=120
Environment=SOCIALBLADE_PAGE_DELAY_SECONDS=20
Environment=SOCIALBLADE_FIREFOX_PROFILE_DIR=/home/orangepi/.local/share/skimmer/socialblade-firefox-profile
Environment=SOCIALBLADE_CLOUDFLARE_BACKOFF_SECONDS=21600
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/orangepi/code_projects/skimmer/src
Environment=SKIMMER_YOUTUBE_CPU=0
ExecStart=/usr/bin/xvfb-run -a -s "-screen 0 1920x1080x24" /usr/bin/python3 -m skimmer.services.workflow
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
